# deploy/systemd/epica.service
[Unit]
Description=Epica Django app (gunicorn)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/epica
EnvironmentFile=/opt/epica/.env
ExecStart=/opt/epica/.venv/bin/gunicorn epica.wsgi:application --bind 127.0.0.1:8001 --workers 3 --timeout 60 --access-logfile - --error-logfile -
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
